# .docker/rootfs/etc/supervisor.d/supervisord.ini
[supervisord]
nodaemon = true

[program:pdns-recursor]
command=pdns_recursor --daemon=no
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.info.log
stderr_logfile=/var/log/supervisor/%(program_name)s.error.log
stdout_logfile_maxbytes=16MB
stderr_logfile_maxbytes=16MB

[program:dnscrypt-wrapper]
autostart={{DNSCRYPT_ON}}
command=dnscrypt-wrapper --resolver-address=127.0.0.1:53 --listen-address=0.0.0.0:{{DNSCRYPT_PORT}} --provider-name={{DNSCRYPT_PROVIDER}} --crypt-secretkey-file=server.key --provider-cert-file=server.cert
directory=/srv/dnscrypt-wrapper
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.info.log
stderr_logfile=/var/log/supervisor/%(program_name)s.error.log
stdout_logfile_maxbytes=16MB
stderr_logfile_maxbytes=16MB

[program:collectd]
autostart={{COLLECTD_ON}}
command=collectd -f
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.info.log
stderr_logfile=/var/log/supervisor/%(program_name)s.error.log
stdout_logfile_maxbytes=16MB
stderr_logfile_maxbytes=16MB

[program:fail2ban]
autostart={{FAIL2BAN_ON}}
command=/usr/bin/fail2ban-server -f
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.info.log
stderr_logfile=/var/log/supervisor/%(program_name)s.error.log
stdout_logfile_maxbytes=16MB
stderr_logfile_maxbytes=16MB
