# .docker/rootfs/etc/supervisor.d/supervisord.ini
[supervisord]
nodaemon = true

[program:named]
command=named -u named -f
redirect_stderr=true
stdout_logfile=/var/log/supervisor/named.info.log
stderr_logfile=/var/log/supervisor/named.error.log
stdout_logfile_maxbytes=16MB
stderr_logfile_maxbytes=16MB

[program:dnscrypt-wrapper]
command=dnscrypt-wrapper --resolver-address=127.0.0.1:53 --listen-address=0.0.0.0:{{DNSCRYPT_PORT}} --provider-name={{DNSCRYPT_PROVIDER}} --crypt-secretkey-file=server.key --provider-cert-file=server.cert
directory=/srv/dnscrypt-wrapper
redirect_stderr=true
stdout_logfile=/var/log/supervisor/dnscrypt-wrapper.info.log
stderr_logfile=/var/log/supervisor/dnscrypt-wrapper.error.log
stdout_logfile_maxbytes=16MB
stderr_logfile_maxbytes=16MB

[program:collectd]
autostart={{COLLECTD_ON}}
command=collectd -f
redirect_stderr=true
stdout_logfile=/var/log/supervisor/collectd.info.log
stderr_logfile=/var/log/supervisor/collectd.error.log
stdout_logfile_maxbytes=16MB
stderr_logfile_maxbytes=16MB
